# common/podman-firewalld-monitor/etc/systemd/system/podman-firewalld-monitor.service
[Unit]
Description=Monitor firewalld reloads and reload Podman networks
After=firewalld.service
Wants=firewalld.service

[Service]
Type=simple
ExecStart=/usr/local/bin/podman-firewalld-monitor.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
